# systemd/mnt-truenas\x2dtailnet-brents\x2ddata.automount
[Unit]
Description=Automount TrueNAS SMB over Tailscale

[Automount]
Where=/mnt/truenas-tailnet/brents-data
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target